# .config/systemd/user/daily-todo-pull.timer
[Unit]
Description=Pull latest code into daily-todos project

[Timer]
# This is in UTC
OnCalendar=*-*-* *:50:00

[Install]
WantedBy=timers.target
