# .config/systemd/user/daily-todo-pull.service
[Unit]
Description=Pull latest code into daily-todos project

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cd ~/projects/daily-todos.sh && git pull'
