# .config/systemd/user/ergodox-watch.service
[Unit]
Description=ErgoDox Watch

[Service]
Type=simple
ExecStart="%h/bin/ergodox-watch"

[Install]
WantedBy=graphical-session.target
