# .config/systemd/user/exocortex-daily-routine.timer
[Unit]
Description=Setup daily routine notes for the coming week

[Timer]
OnCalendar=*-*-* 8:00:00 America/Denver

[Install]
WantedBy=timers.target
