# .config/systemd/user/exocortex-daily-routine.service
[Unit]
Description=Prepare daily routine notes for the upcoming week

[Service]
Type=oneshot
ExecStart=/bin/sh -c '~/projects/exocortex/personal/bin/setup-notes-week'
