# .config/systemd/user/git-autocommit.timer
[Unit]
Description=Commit git repositories

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
