# .config/systemd/user/learn-timers.timer
[Unit]
Description=Learn how systemd user timers work

[Timer]
OnBootSec=15m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
