# .config/systemd/user/org-daily-archive.service
[Unit]
Description=Archive daily todo list and generate future ones

[Service]
Type=oneshot
ExecStart=/bin/sh -c '~/projects/org/bin/daily-archive.sh'
