# .config/systemd/user/org-daily-archive.timer
[Unit]
Description=Archive daily org mode todos

[Timer]
OnCalendar=*-*-* 01:00:00

[Install]
WantedBy=timers.target
