# .config/autostart/copyq.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
Exec=/bin/copyq
Hidden=true
NoDisplay=false
X-GNOME-Autostart-enabled=false
