# .config/overlay/etc/systemd/system/cups.service.d/override.conf
[Service]
RuntimeMaxSec=3600s
Restart=always
